# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'

[program:php-fpm]
command=service php7.4-fpm start --nodaemonize